# .config/systemd/user/sway-session.target
[Unit]
Description=sway session
Documentation=man:systemd.special(7)
BindsTo=graphical-session.target
Wants=graphical-session.target
After=graphical-session.target